# pbs.ini
[supervisord]
nodaemon=true
user=root
loglevel=debug

[unix_http_server]
file = /tmp/supervisord.sock

[program:proxmox-backup-proxy]
command=/usr/lib/x86_64-linux-gnu/proxmox-backup/proxmox-backup-proxy
user=backup
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[program:proxmox-backup-api]
command=/usr/lib/x86_64-linux-gnu/proxmox-backup/proxmox-backup-api
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0